# modules/global/templates/dev.emacs.plist.erb
  <?xml version="1.0" encoding="UTF-8"?>
  <!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
  <plist version="1.0">
    <dict>
      <key>Label</key>
      <string>dev.emacs</string>
      <key>ProgramArguments</key>
      <array>
        <string><%= scope.lookupvar("::boxen::config::homebrewdir") +
                    "/bin/emacs" %></string>
        <string>--daemon=server</string>
      </array>
      <key>KeepAlive</key>
      <true/>
      <key>RunAtLoad</key>
      <true/>
      <key>UserName</key>
      <string><%= scope.lookupvar("::boxen_user") %></string>
    </dict>
  </plist>
  
